# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:fastapi]
command=/app/venv/bin/python -m uvicorn fastapi_app.main:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi_err.log

[program:flask]
command=/app/venv/bin/python flask_app/app.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/flask.log
stderr_logfile=/var/log/supervisor/flask_err.log

[program:laravel]
command=php -S 0.0.0.0:8000 -t /app/laravel_app/public
directory=/app/laravel_app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/laravel.log
stderr_logfile=/var/log/supervisor/laravel_err.log